# configs/multi_protocol_gateway.service
[Unit]
After=gateway.service

[Service]
ExecStart=/usr/local/bin/multi_protocol_gateway.sh

[Install]
WantedBy=default.target
